# Raspbian Scripts/enamel_usb_autostart.rules
ACTION=="add", KERNEL=="sd?1", RUN+="/usr/bin/sudo -u pi /home/pi/startSCALP.sh"